# deploy/pennywise.service
[Unit]
Description=Pennywise Personal Finance Dashboard
After=network.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=pennywise
Group=pennywise
WorkingDirectory=/opt/pennywise
ExecStart=/opt/pennywise/server
Restart=on-failure
RestartSec=10

Environment=PENNYWISE_DB_PATH=/opt/pennywise/data/pennywise.db
Environment=PENNYWISE_PORT=8080
EnvironmentFile=-/opt/pennywise/.env

ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/pennywise/data
NoNewPrivileges=true
PrivateTmp=true

StandardOutput=journal
StandardError=journal
SyslogIdentifier=pennywise

[Install]
WantedBy=multi-user.target
